# sv/zookeeper/run
#!/bin/sh
exec 2>&1
exec /zookeeper/bin/zkServer.sh start-foreground
