# sv/overlord/run
#!/bin/sh

sv start coordinator || exit 1

cd /druid-services
exec 2>&1
exec java -Xmx2g -Duser.timezone=UTC -Dfile.encoding=UTF-8 -classpath lib/*:config/overlord io.druid.cli.Main server overlord
